# etc/crontab
* *     * * *   root    node /opt/aqmonitor/ReefMonitor/scripts/main_job.js